# dpkg/lib/systemd/system/dnsproxy.service
[Install]
WantedBy=multi-user.target

[Unit]
EnvironmentFile=/etc/default/dnsproxy
ExecStart=/usr/bin/dnsproxy ${OPTS}